# rc.local
!/bin/sh -e
#ln -s /data/smartmeter/service /service/smartmeter
ln -s /data/pvinverter/service /service/pvinverter
